# raposa_launchers/services/startup_service.service
[Unit]
After=raposa_robot_start_service

[Service]
ExecStart=/usr/bin/start_raposa

[Install]
WantedBy=default.target